# crontab
0 3 */3 * * /usr/local/bin/python3 /specs_micro/laptop_specs_scraper.py

0 2 */1 * * /usr/local/bin/python3 /specs_micro/update_price.py
